# systemd/.config/systemd/user/backup.service
[Unit]
Description=Backs up files

[Service]
ExecStart=bash -c 'backup.sh ; echo "backuped at `date`" >/home/parth/backup'
